# local/service/screensaver-watch/run
#!/bin/bash
cd "${HOME}"
PARENT=xscreensaver

exec xscreensaver-command -w

