# local/src/service/nm-applet/run
#!/bin/bash
cd "${HOME}"
EXPORT PARENT=nm-applet
until pidof wmsystemtray >> /dev/null 2>> /dev/null;do
	sleep 0.2s
	echo -n .
done
echo '!'

exec /usr/bin/nm-applet
sleep 5m
